# scripts/firstboot/klipos-setup.service
[Unit]
Description=KlipOS First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/etc/klipos/.setup_complete

[Service]
Type=oneshot
ExecStart=/usr/local/bin/klipos-setup.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
